# ok_base-files/etc/hotplug.d/iface/61-hostroute
[ "$INTERFACE" = "lan1" ] && {
    [ "$ACTION" = "ifup" -o "$ACTION" = "ifupdate" ] && {
        {
        sleep 3
        . /lib/functions/network.sh; network_get_gateway _gateway lan1
        route add -host log.networkworld365.com gw $_gateway br-lan1
        } &
    }
}

