# ok_base-files/etc/hotplug.d/iface/63-capwapc
[ "$INTERFACE" = "lan1" ] && {
    [ "$ACTION" = "ifupdate" ] && {
        /etc/init.d/capwapc restart
    }
}

