# ok_base-files/etc/crontabs/root
*/10 * * * * /lib/okos/apstats.sh
*/10 * * * * /lib/okos/setgre.sh
