# ok_base-files/etc/crontabs/root
* * * * * nice -n -20 /lib/okos/runtimefixup.sh
*/1 * * * * /lib/okos/apstats.sh
*/10 * * * * /lib/okos/setgre.sh
