# ok_base-files/etc/crontabs/root
*/1 * * * * export has_ifjson=1;/lib/okos/devstats.sh
