# ok_base-files/etc/crontabs/root
*/10 * * * * /lib/okos/gwstats.sh
*/1 * * * * export has_ifjson=1;/lib/okos/devstats.sh
