# ok_base-files/etc/crontabs/root
*/1 * * * * nice -n -15 /lib/okos/apstats.sh
