# ok_base-files/etc/crontabs/root
* * * * * sleep 5;/lib/okos/runtimefixup.sh
*/10 * * * * /lib/okos/apstats.sh
*/10 * * * * /lib/okos/setgre.sh
