# ok_base-files/etc/init.d/apfw.dyn
#!/bin/sh /etc/rc.common

START=45


. /lib/okos/apfw.dyn

start ()
{
    apfw_dyn_start
}


stop ()
{
    apfw_dyn_stop
}

