# ok_base-files/etc/init.d/apfw
#!/bin/sh /etc/rc.common

START=45


. /lib/okos/apfw.user

start ()
{
    [ ! -f "/tmp/.init_iptables" ] && {

        set_global_static_chains_rules

       touch /tmp/.init_iptables
    }
}

