# ok_base-files/etc/init.d/cmdagent
#!/bin/sh /etc/rc.common

START=90

start ()
{
    nice -n +20 /lib/okos/agent.py start 2>&1 | logger -t cmdagent -p 6
}

stop ()
{
    /lib/okos/agent.py stop 2>&1 | logger -t cmdagent -p 6
}
