# ok_base-files/etc/init.d/predownload
#!/bin/sh /etc/rc.common

START=99


stop ()
{
    pid=`pgrep -f "wget -c $imgurl"`
    [ -n "$pid" ] &&  kill $pid
}


start ()
{
    [ -f "/overlay/root/imgurl" ] && { 
        imgurl="`cat /overlay/root/imgurl`"
        (sleep 1800;ulimit -f 16384;nice -n 20 wget -c "$imgurl" -O "/overlay/root/${imgurl##*/}" 2>&1 | logger -t 'predownload';ulimit -f unlimited;) & 
    }
}

